# ops/systemd/jualan-api.service
[Unit]
Description=Jualan Payment Listener API Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONPATH=__PROJECT_DIR__/src
EnvironmentFile=__PROJECT_DIR__/.env
ExecStart=__PYTHON_BIN__ -m app.run_api
Restart=on-failure
RestartSec=5
User=__RUN_USER__

[Install]
WantedBy=multi-user.target
